# systemd/user/metasphere-reaper.timer
[Unit]
Description=Metasphere reaper timer (every 60s)
Documentation=https://github.com/julianfleck/metasphere-agents/blob/main/docs/OPS.md

[Timer]
OnBootSec=60s
OnUnitActiveSec=60s
AccuracySec=1s
Unit=metasphere-reaper.service

[Install]
WantedBy=timers.target
